# pn_default_xpp_config/etc/profile.d/add_user_bin.sh
# add $HOME/bin to PATH if not present
HOMEBIN="$HOME/bin"
if [[ ! $PATH =~ "$HOMEBIN" ]]; then
	export PATH=$PATH:$HOMEBIN
fi

