# pn_tweak_suspend_on_cover_close/debian/pn-suspend-on-cover.service
[Unit]
Description=

[Service]
ExecStart=/usr/bin/pinenote_sleep_on_cover_close.sh

[Install]
WantedBy=multi-user.target
